# deploy/stockbot-web.service
[Unit]
Description=Random bot dashboard
After=network-online.target

[Service]
User=marvin
WorkingDirectory=/home/marvin/Documents/Stockbot
# No STOCKBOT_WEB_DB: the web app reads .env and auto-selects
# data/dashboard-{paper,live}.db, matching whichever account the bot uses.
ExecStart=/home/marvin/Documents/Stockbot/.venv/bin/python scripts/serve_web.py --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
