# deploy/stockbot-random.service
[Unit]
Description=Random trading bot
After=network-online.target
Wants=network-online.target

[Service]
User=marvin
WorkingDirectory=/home/marvin/Documents/Stockbot
ExecStart=/home/marvin/Documents/Stockbot/.venv/bin/python scripts/random_portfolio.py \
    --execute --record-db auto \
    --label "Random Bot" --poll-interval 30 --trade-prob 0.2
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
